# packaging/coordinode.service
[Unit]
Description=CoordiNode database server
Documentation=https://github.com/structured-world/coordinode
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=coordinode
Group=coordinode
StateDirectory=coordinode coordinode/data
StateDirectoryMode=0750
WorkingDirectory=/var/lib/coordinode
# All server tunables live in the YAML config file (the single config gate).
# Override individual values without editing the file by appending CLI flags
# here via a drop-in (`systemctl edit coordinode`): the command line wins.
ExecStart=/usr/bin/coordinode serve --config /etc/coordinode/coordinode.conf
# Logging is controlled through the environment (standard Rust ecosystem knobs).
# Override via `systemctl edit coordinode`:
#   [Service]
#   Environment=RUST_LOG=debug
#   Environment=COORDINODE_LOG_FORMAT=json
Environment=RUST_LOG=info
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

# Sandboxing
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
ReadWritePaths=/var/lib/coordinode

[Install]
WantedBy=multi-user.target
